# files/services/systemd/gzgspd.service
[Unit]
Description=Guangzhou College of Technology and Business portal login daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/gzgspd --daemon --config /etc/gzgspd/config.json
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
